# cip-opcua-gateway/systemd/cip-opcua-gateway.container
# /etc/containers/systemd/cip-opcua-gateway.container
# Podman Quadlet — drop this file and run: systemctl daemon-reload
# Then: systemctl enable --now cip-opcua-gateway
#
# Requires: podman >= 4.4, RHEL 9 / Fedora 37+

[Unit]
Description=CIP to OPC-UA Gateway
After=network-online.target
Wants=network-online.target

[Container]
# Image built locally — or push to an internal registry and reference it here
Image=localhost/cip-opcua-gateway:latest

# Mount your tags config (read-only into container)
Volume=/etc/cip-opcua-gateway/tags.json:/config/tags.json:ro,z

# OPC-UA port — bind to specific IP if needed: 192.168.1.5:4840:4840
PublishPort=4840:4840

# Environment overrides
Environment=LOG_LEVEL=INFO

# Run as non-root; matches the uid baked into the image
User=1001
Group=1001

# Network access to PLC subnet (adjust or use host if needed)
# Network=host   ← uncomment if PLC is on same L2 as host NIC

# Resource limits
Memory=256m
CPUQuota=50%

# Auto-restart on failure
Restart=on-failure
RestartSec=10s

[Service]
Restart=on-failure
TimeoutStartSec=60

[Install]
WantedBy=default.target
