# cip-opcua-gateway/systemd/cip-opcua-gateway.container
# /etc/containers/systemd/cip-opcua-gateway.container
# Podman Quadlet — drop this file and run: systemctl daemon-reload
# Then: systemctl enable --now cip-opcua-gateway
#
# Requires: podman >= 4.4, RHEL 9 / Fedora 37+

[Unit]
Description=CIP to OPC-UA Gateway
After=network-online.target
Wants=network-online.target

[Container]
# Image built locally — or push to an internal registry and reference it here
Image=localhost/cip-opcua-gateway:latest

# Mount config dir (writable so Web UI can save changes)
Volume=/etc/cip-opcua-gateway/tags.json:/config/tags.json:z

# Use host networking so the container shares the host's real IPs
# (needed for correct OPC-UA endpoint advertisement and direct PLC reach)
Network=host

# Environment overrides
Environment=LOG_LEVEL=INFO

# Run as non-root; matches the uid baked into the image
User=1500
Group=1500

# Network access to PLC subnet (adjust or use host if needed)
# Network=host   ← uncomment if PLC is on same L2 as host NIC

[Service]
Restart=on-failure
RestartSec=10s
TimeoutStartSec=60
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=default.target
